# deployment/argos-final.service
[Unit]
Description=Argos SDR & Network Analysis Console
After=network.target gpsd.service
Wants=gpsd.service

[Service]
Type=simple
User=__SETUP_USER__
Group=__SETUP_USER__
WorkingDirectory=__PROJECT_DIR__
EnvironmentFile=__PROJECT_DIR__/.env
Environment="NODE_ENV=production"
Environment="NODE_OPTIONS=--max-old-space-size=1024"
ExecStart=/usr/bin/node build
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=argos
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
