# deployment/argos-final.service
[Unit]
Description=Argos SDR & Network Analysis Console
After=network.target gpsd.service
Wants=gpsd.service

[Service]
Type=simple
User=kali
Group=kali
WorkingDirectory=/home/kali/Documents/Argos/Argos
EnvironmentFile=/home/kali/Documents/Argos/Argos/.env
Environment="NODE_ENV=production"
Environment="NODE_OPTIONS=--max-old-space-size=1024"
ExecStart=/usr/bin/node build
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=argos
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
