# deployment/argos-dev.service
[Unit]
Description=Argos Development Server
After=network.target

[Service]
Type=simple
User=kali
Group=kali
WorkingDirectory=/home/kali/Documents/Argos/Argos
EnvironmentFile=/home/kali/Documents/Argos/Argos/.env
Environment="NODE_ENV=development"
Environment="NODE_OPTIONS=--max-old-space-size=1024"
ExecStart=/usr/bin/npm run dev
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=argos-dev

[Install]
WantedBy=multi-user.target
