# deployment/argos-process-manager.service
[Unit]
Description=Argos Process Manager
After=network.target

[Service]
Type=simple
User=kali
Group=kali
WorkingDirectory=/home/kali/Documents/Argos/Argos
ExecStart=/usr/local/bin/argos-process-manager.sh monitor
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=argos-process-manager
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
