# deployment/argos-headless.service
[Unit]
Description=Argos Headless Debug Service (Xvfb + Chromium)
After=network.target
Documentation=https://github.com/Graveside2022/Argos

[Service]
Type=simple
User=kali
Environment=DISPLAY=:99
ExecStart=/bin/bash /home/kali/Documents/Argos/Argos/scripts/dev/debug-headless.sh http://localhost:5173
Restart=always
RestartSec=10
# Chromium may need time to shut down gracefully
TimeoutStopSec=15

[Install]
WantedBy=multi-user.target
