# deployment/argos-wifi-resilience.service
[Unit]
Description=Argos WiFi Resilience Manager
After=network.target wpa_supplicant.service NetworkManager.service
Wants=network.target

[Service]
Type=simple
User=kali
Group=kali
WorkingDirectory=/home/kali/Documents/Argos/Argos
ExecStart=/usr/local/bin/argos-wifi-resilience.sh monitor
Restart=on-failure
RestartSec=30
StandardOutput=journal
StandardError=journal
SyslogIdentifier=argos-wifi-resilience
NoNewPrivileges=true
PrivateTmp=true
PrivateNetwork=false

[Install]
WantedBy=multi-user.target
