# deployment/argos-cpu-protector.service
[Unit]
Description=Argos CPU Protector
After=network.target

[Service]
Type=simple
User=kali
Group=kali
WorkingDirectory=/home/kali/Documents/Argos/Argos
ExecStart=/usr/local/bin/argos-cpu-protector.sh monitor
Restart=on-failure
RestartSec=15
StandardOutput=journal
StandardError=journal
SyslogIdentifier=argos-cpu-protector
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
