# deployment/argos-dev-monitor.service
[Unit]
Description=Argos Development Monitor (Vite/Chromium Keepalive)
After=network.target

[Service]
Type=simple
ExecStart=/bin/bash %h/Documents/Argos/Argos/scripts/ops/keepalive-dev.sh
WorkingDirectory=%h/Documents/Argos/Argos
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
