# deployment/argos-droneid.service
[Unit]
Description=Argos DroneID Backend (RemoteIDReceiver)
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=__DRONEID_DIR__/Receiver
ExecStartPre=/bin/sleep 10
ExecStart=__DRONEID_DIR__/Receiver/venv/bin/python3 __DRONEID_DIR__/Receiver/backend/dronesniffer/main.py -p 80
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=argos-droneid
Environment="PYTHONUNBUFFERED=1"

[Install]
WantedBy=multi-user.target
