# deployment/gsmevil-patch.service
[Unit]
Description=GSM Evil Socket.IO Patch Service
Before=gsmevil.service

[Service]
Type=oneshot
ExecStart=/home/kali/Documents/Argos/Argos/scripts/patch-gsmevil-socketio.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
