# deployment/wardragon-fpv-detect.service.d/argos-conflicts.conf
[Unit]
# Installed by Argos — prevents B205/PSU contention. Systemd stops the
# competing unit before this one starts. See plan at
# ~/.claude/plans/just-fix-it-as-typed-squirrel.md
Conflicts=sparrow-wifi-agent.service
