# config/systemd/coral-worker.service
[Unit]
Description=Coral TPU Worker for RSSI Localization
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/projects/Argos
Environment="PATH=/home/ubuntu/projects/Argos/.coral_env/bin"
ExecStart=/home/ubuntu/projects/Argos/.coral_env/bin/python /home/ubuntu/projects/Argos/src/lib/services/localization/coral/coral_worker.py /home/ubuntu/projects/Argos/models/rssi_predictor.tflite
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
